# templates/github-actions-runner.service.epp
<%- | String $instance_name,
      Stdlib::Absolutepath $root_dir,
      String $user,
      String $group,
| -%>
[Unit]
Description=github-actions-runner-<%= $instance_name %>
After=network.target

[Service]
ExecStart=<%= $root_dir %>/<%= $instance_name %>/runsvc.sh
WorkingDirectory=<%= $root_dir %>/<%= $instance_name %>
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=5min
User=<%= $user %>
Group=<%= $group %>

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
