# templates/github-actions-runner.service.epp
[Unit]
Description=github-actions-runner
After=network.target

[Service]
ExecStart=<%= $github_actions_runner::root_dir %>/<%= $github_actions_runner::instance::name %>/runsvc.sh
WorkingDirectory=<%= $github_actions_runner::root_dir %>/<%= $github_actions_runner::instance::name %>
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=5min
User=<%= $github_actions_runner::instance::user %>
Group=<%= $github_actions_runner::instance::group %>

[Install]
WantedBy=multi-user.target
